# deploy/dayclaim-watchdog.timer
[Unit]
Description=Run DayClaim watchdog every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=10s

[Install]
WantedBy=timers.target
